# Embedded_linux_task11/system/task.target
[Unit]
Description= Qt app Target
Requires=mnt-dummy.mount prep.service qt-app-test.service

[Install]
WantedBy=graphical.target